# scripts/systemd/libyaintel-alerts.service
[Unit]
Description=LibyaIntel Alerts Worker
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=simple
User=libyaintel
Group=libyaintel
WorkingDirectory=/opt/libyaintel
EnvironmentFile=/etc/libyaintel/libyaintel.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/libyaintel/.venv/bin/python -u -m backend.runner.alerts.deliver

Restart=always
RestartSec=5

TimeoutStopSec=30
KillSignal=SIGTERM

NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
RemoveIPC=true
CapabilityBoundingSet=
AmbientCapabilities=
ProtectSystem=strict
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectClock=true
ReadWritePaths=/var/lib/libyaintel
StateDirectory=libyaintel
RuntimeDirectory=libyaintel
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ProtectProc=invisible
ProcSubset=pid
PrivateUsers=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal
SyslogIdentifier=libyaintel-alerts
UMask=027

[Install]
WantedBy=multi-user.target
