# systemd/libyaintel-market-quotes.timer
[Unit]
Description=Schedule market quotes fetch

[Timer]
OnBootSec=5m
OnUnitActiveSec=30m
Persistent=true

[Install]
WantedBy=timers.target

