# systemd/libyaintel-gdelt-ingest.timer
[Unit]
Description=Run LibyaIntel GDELT ingest every 30 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
RandomizedDelaySec=120
Persistent=true
Unit=libyaintel-gdelt-ingest.service

[Install]
WantedBy=timers.target
